# scripts/zkaccess-sync/pointaflex-sync-cit.service
[Unit]
Description=PointaFlex ZKTeco Sync - Terminal CIT
Documentation=https://github.com/pointaflex
After=network.target

[Service]
Type=simple
User=assyin
WorkingDirectory=/home/assyin/PointaFlex/scripts/zkaccess-sync
ExecStart=/home/assyin/.nvm/versions/node/v18.20.8/bin/node /home/assyin/PointaFlex/scripts/zkaccess-sync/sync-cit.js
Restart=always
RestartSec=10
StandardOutput=append:/home/assyin/PointaFlex/scripts/zkaccess-sync/sync_cit.log
StandardError=append:/home/assyin/PointaFlex/scripts/zkaccess-sync/sync_cit.log
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
